# code/ansible/firewall-env/files/supervisor/superv-web.conf
[program:superv-web]
stopasgroup=true
user=root
autostart=true
autorestart=true
startsecs=10
startretries=10
directory=/opt/firewall/frontend
command=/usr/sbin/nginx -g 'daemon off;'
redirect_stderr=true
stdout_logfile=/data/logs/superv-web/superv-web.log
priority=12
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2